# PDF/PDFReactor/bin/pdfreactor.service
[Unit]
Description=PDFreactor Web Service

[Service]
Type=forking
ExecStart=/opt/PDFreactor/bin/pdfreactorwebservice start
ExecReload=/opt/PDFreactor/bin/pdfreactorwebservice restart
ExecStop=/opt/PDFreactor/bin/pdfreactorwebservice stop

[Install]
Alias=pdfreactorwebservice pdfreactor
WantedBy=multi-user.target